# rootfs/etc/services.d/discordbot-self-radio/run
#!/usr/bin/with-contenv sh
python /usr/src/app/bot.py
